# debian/vyatta-cfg.postrm.in
#!/bin/bash

set -e

if [ "$1" = "purge" ]; then
	systemctl disable vyatta-router.service > /dev/null || exit $?
fi

exit 0
